# deploy/cleanup.service
[Unit]
Description=Cleanup orphaned GitHub runner droplets

[Service]
Type=oneshot
User=webhook
Group=webhook
EnvironmentFile=/etc/github-runners/env
ExecStart=/usr/local/bin/cleanup
