# deploy/webhook.service
[Unit]
Description=GitHub Actions Runner Webhook Listener
After=network.target

[Service]
Type=simple
User=webhook
Group=webhook
EnvironmentFile=/etc/github-runners/env
ExecStart=/usr/local/bin/webhook
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
